# pi/food-label.service
[Unit]
Description=Food Label Printer App
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/food-label-app
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
